# systemd/system/iwd-firewalld-zone.path
[Path]
Unit=iwd-firewalld-zone.target
PathChanged=/var/lib/iwd

[Install]
WantedBy=multi-user.target

