# systemd/system/iwd-firewalld-zone.target
[Unit]
Description=Target for all %p services
StopWhenUnneeded=yes

[Install]
WantedBy=%p.path
