# systemd/system/iwd-firewalld-zone.path
[Path]
Unit=%p.target
PathChanged=/var/lib/iwd

[Install]
WantedBy=multi-user.target

