# iwd-firewalld-zone/systemd/system/iwd-firewalld-zone.target
[Unit]
Description=Target for all iwd-firewalld-zone services
StopWhenUnneeded=yes

[Install]
WantedBy=iwd-firewalld-zone.path
